# aws_lab_setup/roles/wetty/templates/wetty.service.j2
[Unit]
Description="WETTY: Terminal over HTTP and HTTPS. https://github.com/krishnasrinivas/wetty"
After=network.target auditd.service
[Service]
ExecStart=/usr/bin/wetty.js -p {{ wetty_port }}
Restart=always
[Install]
WantedBy=multi-user.target
